# .config/autostart/nvidia.sh.desktop
[Desktop Entry]
Exec=/home/offpics/scripts/nvidia.sh
Icon=dialog-scripts
Name=nvidia.sh
Path=
Type=Application
X-KDE-AutostartScript=true
